# vm/systemd/cloudflared-named.service
[Unit]
Description=Cloudflare named tunnel for Hermes
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/bin/bash -c '/usr/local/bin/cloudflared tunnel run --token "$(cat /etc/hermes/cloudflared-token)"'
Restart=always
RestartSec=10
StandardOutput=append:/var/log/hermes/cloudflared.log
StandardError=append:/var/log/hermes/cloudflared.log

[Install]
WantedBy=multi-user.target
